# config/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:sshd]
command=/usr/sbin/sshd -D

[program:php5-fpm]
command=/usr/sbin/php5-fpm -c /etc/php5/fpm
user=root
autostart=true
autorestart=true
priority=5

[program:nginx]
command=/usr/sbin/nginx
user=root
autostart=true
autorestart=true
priority=10
stdout_events_enabled=true
stderr_events_enabled=true
